# launchd/com.cuassistant.advisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  CUassistant advisor chat service for launchd.

  Runs the advisor chat HTTP server on ADVISOR_PORT (8770). The bind
  host is ADVISOR_HTTP_HOST, default 127.0.0.1; the campus pilot sets it
  to 0.0.0.0 in the gitignored .env, which config.ts loads relative to
  WorkingDirectory below. It is deliberately NOT set here — a bind host
  that reaches the network is a deployment decision, and this file is
  checked in.

  Security note: ADVISOR_PASSWORD MUST come from the environment —
  vault-injected, or the gitignored .env — and NOT from this plist.
  Secrets do not belong in plist files. The service fails closed if
  ADVISOR_PASSWORD is unset: on loopback every login attempt is
  rejected, and on a non-loopback ADVISOR_HTTP_HOST the process refuses
  to start at all rather than serve an unauthenticated door to the LAN
  (assertAdvisorAuthConfig, src/advisor-server.ts). Both are correct —
  do not "fix" either by hardcoding a password here.

  This service consumes the public MCP servers (8766, 8767) over
  loopback and adds no MCP tools of its own. It never talks to 8765
  (the credentialed server). Tool/policy changes to those servers do
  not require restarting this service — but restarting this service
  ends every in-flight advisor conversation (sessions are in-memory).

  Setup (one time):
    1. Confirm NPM_PATH below matches `which npm` on this host.
    2. Ensure ADVISOR_PASSWORD is set in the environment this launchd
       job inherits (vault-injected), not in this file.
    3. cp launchd/com.cuassistant.advisor.plist ~/Library/LaunchAgents/
    4. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.cuassistant.advisor.plist
    5. launchctl list | grep advisor
    6. curl -s -o /dev/null -w "%{http_code}\n" http://127.0.0.1:8770/

  To unload:
    launchctl bootout gui/$(id -u)/com.cuassistant.advisor
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.cuassistant.advisor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>advisor</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/admin/projects/CUassistant</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/admin</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key>
  <string>/Users/admin/Library/Logs/cuassistant.advisor.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/admin/Library/Logs/cuassistant.advisor.err.log</string>
</dict>
</plist>
