# launchd/com.cuassistant.mcp-catalog-http.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  CUassistant GC CURRICULUM MCP server (HTTP) service for launchd.

  Runs the public-DATA curriculum HTTP MCP server, exposing read-only
  Clemson Graphic Communications curriculum tools (catalog years + degree
  plans). The data is bridged from the separate gc_advisor project via its
  query.py CLI, so gc_advisor must be reachable on the same machine/mount.

  AUTH (required): this server is CAMPUS-REACHABLE and authenticates every
  request against MCP_CATALOG_AUTH_TOKEN — its own bearer, distinct from the
  public server's key and from anything on 8765. The key lives in the repo's
  gitignored .env; do NOT write it into this plist. With the key unset the
  process refuses to start rather than serving open.

  BIND: MCP_CATALOG_HTTP_HOST (default 127.0.0.1; set to 0.0.0.0 in .env for
  campus reach). Deliberately NOT MCP_HTTP_HOST, which is shared with the
  credentialed server on 8765 and must stay on loopback.

  NOTE: the MCP HTTP transport performs no Host/Origin validation, so off
  loopback the bearer is the ONLY gate — there is no rebinding protection.

  A containerized NanoClaw agent reaches it via host.docker.internal:8767
  and must now send `Authorization: Bearer <MCP_CATALOG_AUTH_TOKEN>`. The
  credentialed server (8765) stays loopback-only; see
  com.cuassistant.mcp-http.plist and com.cuassistant.mcp-public-http.plist.

  Ports (track in ~/.dev-ports.yaml): mcp_credentialed 8765, mcp_public 8766,
  mcp_curriculum 8767.

  gc_advisor bridge paths default to /Users/admin/projects/gc_advisor. If
  gc_advisor lives elsewhere, add GC_ADVISOR_PYTHON / GC_ADVISOR_QUERY /
  GC_ADVISOR_DB to EnvironmentVariables below.

  Setup (one time):
    1. Edit this file: replace REPO_PATH with the absolute path to your
       local repo, NPM_PATH with the absolute path to your npm binary
       (run `which npm`), and HOME_PATH with $HOME.
    2. cp launchd/com.cuassistant.mcp-catalog-http.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.cuassistant.mcp-catalog-http.plist
    4. (optional) launchctl start com.cuassistant.mcp-catalog-http  # bring up now
    5. tail -f ~/Library/Logs/cuassistant.mcp-catalog.out.log

  The service starts automatically at login and is kept alive by launchd.
  MCP_CATALOG_HTTP_PORT defaults to 8767 if not set in the environment.

  To unload:
    launchctl unload ~/Library/LaunchAgents/com.cuassistant.mcp-catalog-http.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cuassistant.mcp-catalog-http</string>

  <key>ProgramArguments</key>
  <array>
    <string>NPM_PATH</string>
    <string>run</string>
    <string>mcp:catalog:http</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPO_PATH</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>HOME_PATH</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>HOME_PATH/Library/Logs/cuassistant.mcp-catalog.out.log</string>
  <key>StandardErrorPath</key>
  <string>HOME_PATH/Library/Logs/cuassistant.mcp-catalog.err.log</string>
</dict>
</plist>
